# .github/scripts/levl-log-monitor.service
# =============================================================================
# levl-log-monitor.service — Systemd unit untuk log monitor Laravel
#
# Cara install di VPS:
#   1. Copy file ini ke /etc/systemd/system/
#      sudo cp levl-log-monitor.service /etc/systemd/system/
#
#   2. Buat file env — JANGAN commit file ini, isi sesuai VPS:
#      sudo nano /etc/levl-log-monitor.env
#      Isi:
#        TELEGRAM_BOT_TOKEN=123456:ABC-DEF...
#        TELEGRAM_CHAT_ID=-100123456789
#        DEPLOY_PATH=/var/www/levl          ← sama dengan secrets.DEPLOY_PATH di ci-cd.yml
#        COOLDOWN=30
#
#   3. Aktifkan & jalankan:
#      sudo systemctl daemon-reload
#      sudo systemctl enable levl-log-monitor
#      sudo systemctl start levl-log-monitor
#
#   4. Cek status:
#      sudo systemctl status levl-log-monitor
#      sudo journalctl -u levl-log-monitor -f
# =============================================================================

[Unit]
Description=LEVL Laravel Log Error Monitor (Telegram Notifier)
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=www-data
Group=www-data

# File env — berisi token Telegram & DEPLOY_PATH (tidak di-hardcode di sini)
EnvironmentFile=/etc/levl-log-monitor.env

# DEPLOY_PATH diambil dari env file — sinkron dengan secrets.DEPLOY_PATH di ci-cd.yml
ExecStart=/bin/bash ${DEPLOY_PATH}/.github/scripts/log-monitor.sh \
  "${TELEGRAM_BOT_TOKEN}" \
  "${TELEGRAM_CHAT_ID}" \
  "${DEPLOY_PATH}/storage/logs/laravel.log" \
  "LEVL-BE"

# Restart otomatis jika script mati (misal log file di-rotate)
Restart=always
RestartSec=5s

TimeoutStopSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=levl-log-monitor

[Install]
WantedBy=multi-user.target
